# config/systemd/user/borg.service
[Unit]
Description=Borg Backup

[Service]
Type=oneshot
ExecStart=/home/tassilo/bin/backup_data.sh

[Install]
WantedBy=multi-user.target
