# config/systemd/user/borg.timer
[Unit]
Description=Borg Backup Timer

[Timer]
Unit=borg.service
OnCalendar=hourly

[Install]
WantedBy=timers.target
