# config/systemd/user/borg.timer
[Unit]
Description=Borg Backup Timer

[Timer]
Unit=borg.service
OnCalendar=daily

[Install]
WantedBy=timers.target
