# config/systemd/user/borg.timer
[Unit]
Description=Borg Backup Timer

[Timer]
Unit=borg.service
OnBootSec=30min
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
