# config/systemd/user/writing-tracker.service
[Unit]
Description=writing-tracker

[Service]
ExecStart=/home/tassilo/repos/writing-tracker/writing-tracker.sh

[Install]
WantedBy=default.target
