# config/systemd/user/writing-tracker.timer
[Unit]
Description=Writing Tracker Timer

[Timer]
Unit=writing-tracker.service
OnCalendar=daily

[Install]
WantedBy=timers.target
