# config/systemd/user/writing_tracker.timer
[Unit]
Description=Writing Tracker Timer

[Timer]
Unit=writing_tracker.service
OnCalendar=daily

[Install]
WantedBy=timers.target
