# config/autostart/zotero.desktop
[Desktop Entry]
Type=Application
Exec=zotero
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Zotero
Comment=Open Zotero at startup
